# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:bibcat]
command=/bin/bash -c "nohup uwsgi -s 0.0.0.0:5000 -w run:app"
autostart=true
autorestart=true
